# wswd.service
[Unit]
Description=Clean workspace workdir associations

[Service]
Type=oneshot
ExecStart=wswd clean
